# home/Library/LaunchAgents/com.jasonkuhrt.chezmoi-symlink-heal.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jasonkuhrt.chezmoi-symlink-heal</string>
    <key>ProgramArguments</key>
    <array>
        <string>{{ .chezmoi.homeDir }}/.bun/bin/bun</string>
        <string>{{ .chezmoi.homeDir }}/projects/jasonkuhrt/dotfiles/packages/dotctl/src/bin/dotctl.ts</string>
        <string>heal</string>
        <string>--background</string>
    </array>
    <key>WorkingDirectory</key>
    <string>{{ .chezmoi.homeDir }}/projects/jasonkuhrt/dotfiles</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>StartInterval</key>
    <integer>300</integer>
    <key>StandardOutPath</key>
    <string>{{ .chezmoi.homeDir }}/.local/state/dotfiles-symlink/logs/heal-launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>{{ .chezmoi.homeDir }}/.local/state/dotfiles-symlink/logs/heal-launchd.err.log</string>
</dict>
</plist>
